# aurora.service
[Service]
User=aurora
ExecStart=/bin/bash -lc '/home/aurora/current/run'
Environment=RACK_ENV=production
